# roles/openclaw_acp/templates/openclaw-acp.service.j2
[Unit]
Description=AlphaClaw ACP Gateway
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/alphaclaw start
Restart=on-failure
RestartSec=5
User=alphaclaw
Group=alphaclaw

[Install]
WantedBy=multi-user.target
